# systemd/kdeck-python-goal-runner.service
[Unit]
Description=kdeck Python Goal Runner
After=network-online.target kdeck-api.service

[Service]
Type=simple
WorkingDirectory=%h/work/kdeck
Environment=HOME=%h
Environment=PATH=%h/.nvm/versions/node/v24.16.0/bin:%h/.npm-global/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/work/kdeck/scripts/run_python_goal_runner.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
